# root/etc/systemd/system/proc-acpi.service
[Service]
ExecStart=/usr/local/bin/config-wakeup.sh

[Install]
WantedBy=default.target
